# etc/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
command=/usr/local/sbin/dnsmasq
autostart=true
autorestart=true
user=root
log_stderr=true
logfile=/var/log/dnsmasq.log
